# cron/crontab
0 0 * * * curl -X POST http://notificaciones:8003/resumen-diario
